# scripts/linux/mia@.service
[Unit]
Description=M.I.A - Multimodal Intelligent Assistant
Documentation=https://github.com/yourusername/M.I.A_Multimodal_Inteligent_Assistant
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/M.I.A_Multimodal_Inteligent_Assistant
Environment="PATH=/home/%i/M.I.A_Multimodal_Inteligent_Assistant/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/home/%i/M.I.A_Multimodal_Inteligent_Assistant/.env
ExecStart=/home/%i/M.I.A_Multimodal_Inteligent_Assistant/venv/bin/python -m mia.main --text-only
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/%i/M.I.A_Multimodal_Inteligent_Assistant/logs/mia.log
StandardError=append:/home/%i/M.I.A_Multimodal_Inteligent_Assistant/logs/mia-error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/M.I.A_Multimodal_Inteligent_Assistant/logs
ReadWritePaths=/home/%i/M.I.A_Multimodal_Inteligent_Assistant/memory
ReadWritePaths=/home/%i/M.I.A_Multimodal_Inteligent_Assistant/cache
PrivateTmp=true

[Install]
WantedBy=multi-user.target
